# .config/systemd/user/focus-last.service
[Unit]
Description=focus-last(i3 alt tab)

[Service]
Environment=DISPLAY=:0
ExecStart=%h/.i3/scripts/focus-last.py
Restart=always

