# .bin/lib/arch-extra-setup/udev/rules.d/99-keyboard-hotplug.rules
SUBSYSTEM=="input", KERNEL=="input*", ATTRS{phys}=="usb*input*", ATTRS{name}=="*HHKB*", ACTION=="add|remove", RUN+="/usr/local/bin/auto_xmodmap.sh"
