# system-etc/udev/rules.d/95-monitor-hotplug.rules
SUBSYSTEM=="drm", ACTION=="change", ENV{DISPLAY_LAYOUT}="left", ENV{DISPLAY_ROTATE}="normal", RUN+="/usr/local/bin/detect_displays.sh"
